# runner/s6-overlay/s6-rc.d/dockerd/finish
#!/usr/bin/env bash

rm -f /var/run/docker.pid
