# runner/s6-overlay/s6-rc.d/runner/run
#!/usr/bin/env bash

# Cleanup any leftover credentials
rm -f /tmp/private.pem
rm -f /run/runner.token

# run as runner user in a clean environment without secrets
exec su - runner -c "/home/runner/run.sh" 2>&1
