# runner/s6-overlay/s6-rc.d/healthcheck/finish
#!/usr/bin/env bash

if [ "$1" -ne 0 ]; then
  echo "$1" >/run/s6-linux-init-container-results/exitcode
  /run/s6/basedir/bin/halt
fi
