# runner/s6-overlay/s6-rc.d/remco/up
#!/command/execlineb -P
# shellcheck disable=SC1008
remco --config /etc/remco/config.toml
